# pkgs/freetakserver/systemd/freetakserver.service
[Unit]
Description=FreeTAKServer - Open Source TAK Server
Documentation=https://freetakteam.github.io/FreeTAKServer-User-Docs/
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=freetakserver
Group=freetakserver
ExecStart=/opt/conda/bin/python3 -m FreeTAKServer.controllers.services.FTS
WorkingDirectory=/var/lib/freetakserver
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Environment variables
Environment=FTS_FIRST_START=false
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/freetakserver
ReadWritePaths=/var/log/freetakserver
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Network settings - TAK typically uses ports 8080, 8443, 8444, 8446
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freetakserver

[Install]
WantedBy=multi-user.target
Alias=fts.service
